# rootfs/etc/s6-overlay/s6-rc.d/getty/run
#!/command/execlineb -P
exec getty -L -n -l /bin/bash 38400 console linux
