# rootfs/etc/s6-overlay/s6-rc.d/mythbackend/run
#!/command/execlineb -P
exec mythbackend --user mythtv
