# rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/command/with-contenv sh
##!/command/execlineb -P
exec /usr/local/etc/init.sh
